# open-electricity.service
[Unit]
Description=Open Electricity price proxy
After=network.target

[Service]
Type=simple
User=tom
WorkingDirectory=/home/tom/oe-mirror
EnvironmentFile=/home/tom/oe-mirror/.env
ExecStart=/home/tom/oe-mirror/.venv/bin/python /home/tom/oe-mirror/app.py
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
